# rpms/bpac-cmc-user/src/etc/profile.d/bricata_cmc.sh
export SENSOR_SSH_PRIVATE_KEY_DEFAULTS="/home/cmc/.ssh/id_rsa"
